# disable-xhc-wakeup.service
[Unit]
Description=Disable XHC Wakeup
 
[Service]
ExecStart=/usr/local/bin/disable-xhc-wakeup.sh
 
[Install]
WantedBy=multi-user.target
