# .config/autostart/redshift-gtk.desktop
[Desktop Entry]
Name[pt_BR]=Redshift
Comment[cs]=Nástroj pro přizpůsobení barevné teploty
GenericName[fr]=Réglage de la température de couleur
Comment[zh_CN]=色温调节工具
Comment[es]=Herramienta para el ajuste de la temperatura de color
Name[fr]=Redshift
Name[lt]=Redshift
Comment[da]=Justeringsværktøj for farvetemperatur
Hidden=false
StartupNotify=true
GenericName[hu]=Szín hőmérséklet beállítás
Categories=Utility;
GenericName[da]=Justering af farvetemperatur
GenericName[cs]=Přizpůsobení barevné teploty
GenericName[pl]=Dostosowanie temperatury barwowej
Name[hu]=Redshift
Comment[fr]=Outil de réglage de la température de couleur
GenericName[zh_CN]=色温调节
Comment[de]=Farbtemperatur Anpassungswerkzeug
Comment[tr]=Renk sıcaklığı ayarlama aracı
GenericName[nl]=Kleurtemperatuur aanpassing
GenericName[es]=Ajuste de la temperatura de color
Name[de]=Redshift
Name[tr]=Redshift
X-GNOME-Autostart-enabled=true
Name[ru]=Redshift
Name[pl]=Redshift
Name[es]=Redshift
Name[it]=Redshift
Comment[nl]=Kleurtemperatuur aanpassingsprogramma
Comment=Color temperature adjustment tool
Terminal=false
GenericName[lt]=Spalvos temperatūros reguliavimas
Version=1.0
Comment[lt]=Spalvos temperatūros reguliavimo įrankis
GenericName[it]=Regolazione della temperatura del colore
Comment[ru]=Инструмент регулирования цветовой температуры
GenericName[ru]=Настройка цветовой температуры
Exec=redshift-gtk
GenericName[tr]=Renk sıcaklığı ayarı
Name[zh_CN]=红移
Type=Application
GenericName[de]=Farbtemperatur Anpassung
GenericName=Color temperature adjustment
Comment[it]=Strumento per la regolazione della temperatura del colore
Icon=redshift
Name=Redshift
Comment[hu]=Színhőmérséklet beállító eszköz
Comment[pl]=Narzędzie do dostosowywania temperatury barwowej

